# systemd/tennisboss.service
[Unit]
Description=TennisBoss Tennis Prediction Engine
After=network.target

[Service]
Type=simple
User=tennisboss
WorkingDirectory=/opt/tennisboss
Environment="PATH=/opt/tennisboss/.venv/bin"
ExecStart=/opt/tennisboss/.venv/bin/python3 run.py serve --host 0.0.0.0 --port 8000
Restart=always
RestartSec=10s

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/tennisboss/state /opt/tennisboss/logs

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tennisboss

[Install]
WantedBy=multi-user.target
