# systemd/tennisboss-scheduler.service
[Unit]
Description=TennisBoss Scheduler (ingest tennis-data.co.uk 6h + recalibration surface 1h)
After=network.target tennisboss-bot.service

[Service]
Type=simple
User=alchemist
WorkingDirectory=/mnt/c/Users/donpa/TennisBoss
EnvironmentFile=/mnt/c/Users/donpa/TennisBoss/.env
ExecStart=/usr/bin/python3 -c "from bot.scheduler import run_scheduler; run_scheduler()"
Restart=always
RestartSec=10s

StandardOutput=journal
StandardError=journal
SyslogIdentifier=tennisboss-scheduler

[Install]
WantedBy=multi-user.target
